# docker/supervisord.conf
[supervisord]
nodaemon = true
logfile = /var/log/supervisord.conf
pidfile = /var/run/supervisord.pid

[program:nginx]
command = /usr/sbin/nginx -c /app/timtec/docker/nginx.conf
startsecs = 5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:app-gunicorn]
command = /usr/local/bin/gunicorn timtec.wsgi 
    -u timtec 
    -g www-data 
    -b 0.0.0.0:8000 
    --log-level=%(ENV_GUNICORN_LOG_LEVEL)s 
    --log-file=/dev/stdout %(ENV_GUNICORN_EXTRA_FLAGS)s
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
